# backend/deploy/tetramem-tcp.service
[Unit]
Description=TetraMem v14 Memory Server (TCP)
After=network.target

[Service]
Type=simple
User=tetramem
Group=tetramem
WorkingDirectory=/opt/tetramem
ExecStart=/opt/tetramem/run-tcp.sh
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

Environment=TETRAMEM_DATA_DIR=/var/lib/tetramem
Environment=TETRAMEM_MASTER_KEY=/3NwoHqiE5WWZxyls2r3yL7SwHn3GY6Ek26LgHZwXj4=
Environment=TETRAMEM_API_KEY=Q3EOKNmcC92M0iogdbADZnPVSf5eGvr4
Environment=DEEPSEEK_API_KEY=sk-147214015ffa49b6b728c5709e110aec
Environment=RUST_LOG=tetramem_v14=info,tetramem_mcp=info
Environment=EMBEDDING_API_URL=disabled://none
Environment=TETRAMEM_PORT=19100

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tetramem

[Install]
WantedBy=multi-user.target
